# data/conf/mosn.supervisor.conf
[program:mosn]
command=/home/admin/mosn/bin/mosn start -c /home/admin/mosn/conf/mosn_config.json --service-cluster=gateway --service-node=docker-compose-node-2 --sm=role:gloo-system~gateway-proxy
user=admin
startsecs=1
autostart=true
autorestart=true
stdout_events_enabled = true
stderr_events_enabled = true
